# streamer-ui/packages/server/cron/crontab
0 23 * * 0 /opt/streamer-ui-server/purge-db.sh >/dev/null 2>&1