# streamer/supervisord.conf
[supervisord]
nodaemon=true

[program:streamer]
command=/opt/streamer/start_streamer.sh
stdout_logfile=/opt/streamer/log/streamer.out
stderr_logfile=/opt/streamer/log/streamer.err
