# streamer-agent/streamer-agent-opm-1.0/debian/streamer-agent-opm.streamer-agent-opm-cli.user.service
[Unit]
Description=Streamer Agent Service client

[Service]
Type=onshot
RemainAfterExit=yes
ExecStop=/usr/bin/streamer-agent-opm-cli -e

[Install]
WantedBy=multi-user.target
