# streamer-agent/streamer-agent-opm-1.0/debian/streamer-agent-opm.service
[Unit]
Description=Streamer Agent Service for FieldLine OPM

[Service]
EnvironmentFile=/etc/streamer-agent.cfg
ExecStart=/sbin/streamer-agent

[Install]
WantedBy=multi-user.target
