# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:devcron]
command=/opt/python/bin/devcron.py /cron/crontab
stdout_logfile=/opt/meg-dataflow/log/cron.out
stderr_logfile=/opt/meg-dataflow/log/cron.err
